# deploy/supervisor.conf
[program:wuthering-waves-dps]
command=/root/ai/WutheringWavesDPS/backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 14876 --workers 4
directory=/root/ai/WutheringWavesDPS/backend
user=www-data
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/wuthering-waves-dps/backend.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=PYTHONUNBUFFERED="1"
